# doorbell.service

[Unit]
Description=Accessible Doorbell (Flask + Picamera2)
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/accessible-doorbell
ExecStart=/usr/bin/python3 /home/pi/accessible-doorbell/app.py
Restart=on-failure
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
